# docker/supervisor.conf
[program:server]
command = node /app/index.js
user = node
autostart=true
autorestart=true
stdout_logfile = /logs/review-me-out.log
stderr_logfile = /logs/review-me-err.log